# roles/borg/templates/borgmatic.timer.j2
[Unit]
Description=Run borgmatic backup

[Timer]
OnCalendar={{ BORGMATIC_BACKUP_FREQUENCY }}
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
